# systemd/llama-server.service
[Unit]
Description=llama-server (Vulkan - Strix Halo - Qwen 3.6 35B-A3B + native MTP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=GGML_VK_PREFER_HOST_MEMORY=ON
ExecStart=/home/awpapa/llama.cpp/build/bin/llama-server \
  -m /home/awpapa/models/qwen3.6-mtp/Qwen3.6-35B-A3B-UD-Q4_K_XL.gguf \
  --spec-type draft-mtp \
  --spec-draft-n-max 3 \
  --ctx-size 262144 \
  --port 8001 \
  --host 0.0.0.0 \
  --no-warmup \
  -ngl 99 \
  --mmap \
  --jinja \
  --ubatch-size 1024 \
  -fa 1 \
  -ctk q8_0 \
  -ctv q8_0 \
  --cache-prompt \
  --reasoning-budget 4096 \
  --reasoning-format auto
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
