# systemd/llama-server.service
[Unit]
Description=llama-server (Vulkan - Strix Halo - Qwen3.5-122B-A10B)
After=default.target

[Service]
Type=simple
# llama.cpp Vulkan build from ~/llama-cpp-turboquant (turboquant fork).
# Turbo KV cache types are CPU-only; on Vulkan we use q8_0 KV cache.
#
# IMPORTANT flags:
#   --mmap           REQUIRED for Vulkan. Maps weights into 96GB unified VRAM.
#                    Without it, the 72GB model loads into 31GB system RAM → swap death.
#   -fa 1            Flash attention — required for good performance on Strix Halo.
#   -ctk q8_0        KV cache type for keys (turbo2 not supported on Vulkan).
#   -ctv q8_0        KV cache type for values.
#   --ubatch-size 1024  Sweet spot for MoE prefill on Vulkan (tested 512/1024/2048).
#   --cache-prompt   Cache prompt tokens for faster re-processing.
#   --reasoning-budget 500  Cap thinking tokens to avoid runaway reasoning.
#
# Performance (kernel 7.0-rc6, Mesa 25.3.6, RADV):
#   pp: ~393 t/s (2K tokens)
#   tg: ~22 t/s
#   TTFT: ~430ms
ExecStart=/home/awpapa/llama-cpp-turboquant/build/bin/llama-server \
  -m /mnt/Esuna/models/unsloth-qwen3.5-122b-a10b/UD-Q4_K_XL/Qwen3.5-122B-A10B-UD-Q4_K_XL-00001-of-00003.gguf \
  --ctx-size 65536 \
  --port 8001 \
  --host 0.0.0.0 \
  --no-warmup \
  -ngl 99 \
  --mmap \
  --jinja \
  --reasoning-format auto \
  --ubatch-size 1024 \
  -fa 1 \
  -ctk q8_0 \
  -ctv q8_0 \
  --cache-prompt \
  --reasoning-budget 500 \
  --reasoning-budget-message "Thinking budget reached. Stop thinking and answer directly."
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
